# init/uzbl-session-save.timer
[Unit]
Description=uzbl session autosave

[Timer]
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
