# init/uzbl-session@.timer
[Unit]
Description=uzbl session autosave
PartOf=uzbl-manager@.service uzbl-session-save@.service

[Timer]
OnUnitActiveSec=3min
